# crontab
# m h  dom mon dow   command
*/15 * * * * (cd /srv/senderscore-alter ; ./senderscore.py --cron mx1.uni-trier.de)